# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[program:fastapi]
command=/usr/local/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/fastapi.err.log
stdout_logfile=/var/log/supervisor/fastapi.out.log

[program:streamlit]
command=/usr/local/bin/streamlit run viewer_streamlit.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/streamlit.err.log
stdout_logfile=/var/log/supervisor/streamlit.out.log